# crontab
*/5 * * * * php -f /var/www/html/cron.php > /dev/null 2>&1
